# packaging/debian/prerm
#!/bin/sh
set -eu

if [ "${1:-}" = "remove" ] && [ -d /run/systemd/system ]; then
    systemctl disable --now lts-agent.timer || true
fi

exit 0
